# sentry_app/systemd/sentry-agent.service
[Unit]
Description=Smart Sentry fleet agent (BLE + OTA owner)
After=bluetooth.target network-online.target
Wants=bluetooth.target

[Service]
# Runs independent of any graphical/user session so OTA and monitoring
# survive the dashboard restarting or nobody being logged in at all - see
# sentry_agent/agent_main.py's docstring for why the dashboard no longer
# owns BLE directly.
User=klemjo
WorkingDirectory=/home/klemjo/sentry_app/sentry_agent
ExecStart=/home/klemjo/sentry-venv/bin/python3 -u agent_main.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
